# device-net-snmp-el9/src/main/resources/etc/systemd/system/snmpd.service.d/local.conf
[Unit]
After=snmpd-autodiscovery.service
Requisite=snmpd-autodiscovery.service

